# quadlets/alertmanager.container
[Unit]
Description=Alertmanager - Alert Routing and Notification
After=network-online.target monitoring-network.service
Wants=monitoring-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-06-23 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = adopt new digest (pin-container-image.sh --adopt), bake, restart.
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=quay.io/prometheus/alertmanager@sha256:af26fbe4dd1886ac0efd7bd55cd9027da262e105b137a376522b7c14c3626e4a
ContainerName=alertmanager
NoNewPrivileges=true
HostName=alertmanager
Network=systemd-reverse_proxy:ip=10.89.2.72
Network=systemd-monitoring:ip=10.89.4.72

# Configuration and data
# Directory mount (GH#276 pattern): avoids the single-file bind inode trap and
# brings remediation-url.secret (gitignored) in alongside alertmanager.yml.
Volume=%h/containers/config/alertmanager:/etc/alertmanager:ro,Z
Volume=%h/containers/data/alertmanager:/alertmanager:Z

# SMTP password for email notifications
Secret=smtp_password

# Port 9093 NOT published to host - only accessible via:
#   - Monitoring network: 10.89.4.72:9093 (Prometheus, Grafana)
#   - Reverse proxy network: 10.89.2.72:9093 (internal only, no Traefik route)

# Command arguments
Exec=--config.file=/etc/alertmanager/alertmanager.yml --storage.path=/alertmanager --log.level=info --web.external-url=https://alertmanager.patriark.org

# ADR-031 Phase 0 (D2): independent fallback so a Pi-hole outage can't silence alerting.
# 9.9.9.9 is INTERIM — Phase 3 cutover (GH#292, design D-7) replaces these with:
#   DNS=192.168.1.72   (VIP)  /  DNS=192.168.1.69  /  DNS=192.168.1.169
DNS=192.168.1.69
DNS=9.9.9.9

# Health check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:9093/-/healthy || exit 1
HealthInterval=10s
HealthTimeout=3s
HealthRetries=3

# Note: Traefik routing configured in ~/containers/config/traefik/dynamic/routers.yml
# Per design principles: routing belongs in dynamic YAML, not container labels

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=120
# Boot tier A (200): the alerting path must not be starved during the boot
# I/O storm — ADR-035 keeps Alertmanager early and un-deferred so paging works
# while everything else is still thrashing. Promoted explicitly 2026-06-12
# (was implicit default 100, the only alert-path unit without a weight).
StartupCPUWeight=200
StartupIOWeight=200

# Resource Limits
MemoryMax=256M
MemoryHigh=230M

[Install]
WantedBy=default.target
