# quadlets/alert-discord-relay.container
[Unit]
Description=Alert Discord Relay - Alertmanager to Discord Webhook Bridge
After=network-online.target monitoring-network.service
Wants=network-online.target
Requires=monitoring-network.service

[Container]
Image=localhost/alert-discord-relay:latest
ContainerName=alert-discord-relay
HostName=alert-discord-relay
Network=systemd-monitoring

# Discord webhook URL
Secret=discord_webhook_url,type=env,target=DISCORD_WEBHOOK_URL

# Health check (check if service is listening)
# Using python3 since wget is not available in this container
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://localhost:9095/health', timeout=5)" || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=on-failure
TimeoutStartSec=60

# Resource Limits
MemoryMax=128M

[Install]
WantedBy=default.target
