# quadlets/alert-discord-relay.container
[Unit]
Description=Alert Discord Relay - Alertmanager to Discord Webhook Bridge
After=network-online.target reverse_proxy-network.service monitoring-network.service
Wants=monitoring-network.service network-online.target reverse_proxy-network.service

[Container]
# Local build, immutable tag (proton-bridge convention; ADR-030 P5 Tier-2). Tag = build date.
# Build inputs are pinned in config/alert-discord-relay/Dockerfile (FROM @sha256 + hash-locked
# requirements.lock). Recorded build digest: sha256:4df80f75b8100aca7eb699533ac72316759f22848a1b623129e9ed3c0f557e31
# Bump = rebuild, retag with new date, update this line + digest, daemon-reload + restart.
Image=localhost/alert-discord-relay:2026-05-23
ContainerName=alert-discord-relay
NoNewPrivileges=true
HostName=alert-discord-relay
# reverse_proxy first for default route (Discord API access)
Network=systemd-reverse_proxy:ip=10.89.2.86
Network=systemd-monitoring

# ADR-031 Phase 0 (D2): explicit resolvers with an independent fallback (previously inherited
# reverse_proxy DNS=192.168.1.69 only). A Pi-hole outage must not stop alert delivery.
# 9.9.9.9 is INTERIM — replace with the keepalived VIP + both node IPs in Phase 3.
DNS=192.168.1.69
DNS=9.9.9.9

# Discord webhook URL
Secret=discord_webhook_url,type=env,target=DISCORD_WEBHOOK_URL

# Health check (check if service is listening)
# Using python3 since wget is not available in this container
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://localhost:9095/health', timeout=5)" || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=60
# Boot tier C weights (yield disk during fan-out), but UN-deferred: ADR-035 D2
# says the alerting path (Alertmanager + relay) stays early so paging is not
# blinded post-reboot. The After=traefik.service deferral from PR #242 predated
# that decision and contradicted it — removed 2026-06-12.
StartupCPUWeight=50
StartupIOWeight=50

# Resource Limits
MemoryMax=128M
MemoryHigh=115M

[Install]
WantedBy=default.target
