# quadlets/alertmanager.container
[Unit]
Description=Alertmanager - Alert Routing and Notification
After=network-online.target monitoring-network.service
Wants=monitoring-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-23 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=quay.io/prometheus/alertmanager@sha256:51a825c2a40acc3e338fdd00d622e01ec090f72be2b3ea46be0839cd47a4d286
ContainerName=alertmanager
HostName=alertmanager
Network=systemd-reverse_proxy:ip=10.89.2.72
Network=systemd-monitoring:ip=10.89.4.72

# Configuration and data
Volume=%h/containers/config/alertmanager/alertmanager.yml:/etc/alertmanager/alertmanager.yml:ro,Z
Volume=%h/containers/data/alertmanager:/alertmanager:Z

# SMTP password for email notifications
Secret=smtp_password

# Port 9093 NOT published to host - only accessible via:
#   - Monitoring network: 10.89.4.72:9093 (Prometheus, Grafana)
#   - Reverse proxy network: 10.89.2.72:9093 (internal only, no Traefik route)

# Command arguments
Exec=--config.file=/etc/alertmanager/alertmanager.yml --storage.path=/alertmanager --log.level=info --web.external-url=https://alertmanager.patriark.org

# ADR-031 Phase 0 (D2): independent fallback so a Pi-hole outage can't silence alerting.
# 9.9.9.9 is INTERIM — replace with the keepalived VIP + both node IPs in Phase 3.
DNS=192.168.1.69
DNS=9.9.9.9

# Health check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:9093/-/healthy || exit 1
HealthInterval=10s
HealthTimeout=3s
HealthRetries=3

# Note: Traefik routing configured in ~/containers/config/traefik/dynamic/routers.yml
# Per design principles: routing belongs in dynamic YAML, not container labels

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=120
# Boot tier A (200): the alerting path must not be starved during the boot
# I/O storm — ADR-035 keeps Alertmanager early and un-deferred so paging works
# while everything else is still thrashing. Promoted explicitly 2026-06-12
# (was implicit default 100, the only alert-path unit without a weight).
StartupCPUWeight=200
StartupIOWeight=200

# Resource Limits
MemoryMax=256M
MemoryHigh=230M

[Install]
WantedBy=default.target
