# quadlets/audiobookshelf.container
# ~/.config/containers/systemd/audiobookshelf.container
# Audiobookshelf - Self-hosted Audiobook & Podcast Server
# Created: 2026-02-28
# Documentation: https://www.audiobookshelf.org/docs

# ═══════════════════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════════════════
[Unit]
Description=Audiobookshelf Audiobook & Podcast Server
Documentation=https://github.com/advplyr/audiobookshelf
After=network-online.target traefik.service
Wants=network-online.target traefik.service

# ═══════════════════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════════════════
[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 2.35.0, resolved 2026-05-24 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=ghcr.io/advplyr/audiobookshelf@sha256:89276ff2e0b3d2f07dd334b641f27a34ab7f02e1047c60b7b8a30126cb0813a5
ContainerName=audiobookshelf
HostName=audiobookshelf

# Network (single network - standalone service, no inter-service communication)
Network=systemd-reverse_proxy:ip=10.89.2.74

# Environment
Environment=TZ=Europe/Oslo
Environment=PORT=80

# NOTE: ACCESS_TOKEN_EXPIRY=86400 workaround (PR #185, bug #4630) removed in 2.35.0 bump —
# upstream access-token refresh grace period (PR #5004) ships in this release. See issue #188.
# Verify token refresh + Prologue session survival across the access-token boundary post-deploy.

# Config & metadata (exclusive - NOCOW for SQLite)
Volume=/mnt/btrfs-pool/subvol7-containers/audiobookshelf/config:/config:Z
Volume=/mnt/btrfs-pool/subvol7-containers/audiobookshelf/metadata:/metadata:Z

# Media libraries (shared read-only)
Volume=/mnt/btrfs-pool/subvol4-multimedia/Lydbøker:/audiobooks:ro,z
Volume=/mnt/btrfs-pool/subvol4-multimedia/Bøker:/ebooks:ro,z
Volume=/mnt/btrfs-pool/subvol4-multimedia/Podder:/podcasts-archive:ro,z

# Podcast downloads (writable, shared context under multimedia)
Volume=/mnt/btrfs-pool/subvol4-multimedia/audiobookshelf:/podcasts:z

# Health check (no curl in official image, use wget)
HealthCmd=wget --no-verbose --tries=1 --spider http://127.0.0.1:80/healthcheck || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

# ═══════════════════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════════════════
[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
RestartSec=10s

# Resource Limits
MemoryMax=2G
MemoryHigh=1600M
MemorySwapMax=1G

# ═══════════════════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════════════════
[Install]
WantedBy=default.target
