# quadlets/audiobookshelf.container
# ~/.config/containers/systemd/audiobookshelf.container
# Audiobookshelf - Self-hosted Audiobook & Podcast Server
# Created: 2026-02-28
# Documentation: https://www.audiobookshelf.org/docs

# ═══════════════════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════════════════
[Unit]
Description=Audiobookshelf Audiobook & Podcast Server
Documentation=https://github.com/advplyr/audiobookshelf
After=network-online.target restricted-egress-network.service traefik.service
Wants=network-online.target restricted-egress-network.service traefik.service

# ═══════════════════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════════════════
[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 2.35.0, resolved 2026-05-24 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=ghcr.io/advplyr/audiobookshelf@sha256:89276ff2e0b3d2f07dd334b641f27a34ab7f02e1047c60b7b8a30126cb0813a5
ContainerName=audiobookshelf
NoNewPrivileges=true
HostName=audiobookshelf

# Network: restricted-egress (ADR-045, GH#274) — outbound limited to LAN + container
# fabric; INBOUND is unchanged and governed solely by the Traefik router (remote
# access from anywhere still works — Traefik dials in cross-bridge via the pinned
# IP in config/traefik/hosts; replies to 10.89.2.69 are in the filter's accept set).
# Zero internet egress observed in the 18-day baseline (podcasts are archived local
# files: feedURL empty, autoDownloadEpisodes=0 for all 32).
# CEREMONY (accepted 2026-07-17): "Match" metadata lookups and adding a LIVE podcast
# feed need internet → temporarily revert this line to
#   Network=systemd-reverse_proxy:ip=10.89.2.74  (+ hosts file entry back to 10.89.2.74)
# restart, do the fetch, then re-apply this line and restart again.
# Single-homing invariant: restricted-egress must stay the ONLY non-internal network.
Network=systemd-restricted-egress:ip=10.89.11.74

# Environment
Environment=TZ=Europe/Oslo
Environment=PORT=80

# NOTE: ACCESS_TOKEN_EXPIRY=86400 workaround (PR #185, bug #4630) removed in 2.35.0 bump —
# upstream access-token refresh grace period (PR #5004) ships in this release. See issue #188.
# Verify token refresh + Prologue session survival across the access-token boundary post-deploy.

# Config & metadata (exclusive - NOCOW for SQLite)
Volume=/mnt/btrfs-pool/subvol7-containers/audiobookshelf/config:/config:Z
Volume=/mnt/btrfs-pool/subvol7-containers/audiobookshelf/metadata:/metadata:Z

# Media libraries (shared read-only)
Volume=/mnt/btrfs-pool/subvol4-multimedia/Lydbøker:/audiobooks:ro,z
Volume=/mnt/btrfs-pool/subvol4-multimedia/Bøker:/ebooks:ro,z
Volume=/mnt/btrfs-pool/subvol4-multimedia/Podder:/podcasts-archive:ro,z
# qBittorrent 'utdanning' category save path — never had an ABS library folder,
# so its 23 items scanned into stale /ebooks library rows with no matching disk path (isMissing=1)
Volume=/mnt/btrfs-pool/subvol4-multimedia/utdanning:/utdanning:ro,z

# Podcast downloads (writable, shared context under multimedia)
Volume=/mnt/btrfs-pool/subvol4-multimedia/audiobookshelf:/podcasts:z

# Health check (no curl in official image, use wget)
HealthCmd=wget --no-verbose --tries=1 --spider http://127.0.0.1:80/healthcheck || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

# ═══════════════════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════════════════
[Service]
# Tier B (default 100) — ADR-035 boot tier: application class, deliberate default weight
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
RestartSec=10s

# Resource Limits
MemoryMax=2G
MemoryHigh=1600M
MemorySwapMax=1G
CPUQuota=200%

# ═══════════════════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════════════════
[Install]
WantedBy=default.target
