# quadlets/collabora.container
[Unit]
Description=Collabora Online Office Suite
After=network-online.target nextcloud.service
Wants=network-online.target


[Container]
Image=docker.io/collabora/code:latest
ContainerName=collabora
AutoUpdate=registry

# Nextcloud domain whitelist
Environment=aliasgroup1=https://nextcloud.patriark.org

# Collabora configuration
Environment=extra_params=--o:ssl.enable=false --o:ssl.termination=true

# Locale
Environment=dictionaries=en_US no_NO

# Admin credentials (using podman secrets)
Environment=username=admin
Secret=collabora_admin_password,type=env,target=password

# Capabilities for LibreOffice
AddCapability=MKNOD

# Networks (reverse_proxy FIRST for internet access!)
Network=systemd-reverse_proxy:ip=10.89.2.74
Network=systemd-nextcloud:ip=10.89.10.74

# Health check
HealthCmd=curl -f http://localhost:9980/hosting/discovery
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5

# Traefik routing configured via file provider (config/traefik/dynamic/routers.yml)
# No labels needed - cleaner separation of concerns

[Service]
Restart=always
TimeoutStartSec=300

# Resource limits (systemd memory control)
MemoryMax=2G
MemoryHigh=1.8G

[Install]
WantedBy=default.target
