# quadlets/crowdsec.container
[Unit]
Description=CrowdSec Security Engine
After=network-online.target
Wants=network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-23 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — egress-tier, deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=docker.io/crowdsecurity/crowdsec@sha256:2f527c9bb8b367120eb08b82890aa912ce96bfa1ada93dda0721700e4b4e0dde
ContainerName=crowdsec
NoNewPrivileges=true
Network=systemd-reverse_proxy:ip=10.89.2.70

# Volumes
Volume=%h/containers/data/crowdsec/db:/var/lib/crowdsec/data:Z
Volume=%h/containers/data/crowdsec/config:/etc/crowdsec:Z
# Traefik access log acquisition (read-only). Written by traefik.container.
# Shared label (:z) — Traefik writes, CrowdSec tails.
Volume=/mnt/btrfs-pool/subvol7-containers/traefik-logs:/var/log/traefik:ro,z

# Environment - Install Traefik collection
Environment=COLLECTIONS=crowdsecurity/traefik crowdsecurity/http-cve
Environment=GID=1000

# Health check
HealthCmd=cscli version
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
# Tier B (default 100) — ADR-035 boot tier: application class, deliberate default weight
Slice=container.slice
Restart=on-failure
TimeoutStartSec=900

# Resource Limits
MemoryMax=512M
MemoryHigh=460M
CPUQuota=200%

[Install]
WantedBy=default.target
