# quadlets/forgejo.container
[Unit]
Description=Forgejo Git Forge
After=network-online.target forgejo-db.service reverse_proxy-network.service forgejo-network.service
Wants=network-online.target
Requires=forgejo-db.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 15, resolved 2026-05-23 (index digest, multi-arch)
Image=codeberg.org/forgejo/forgejo@sha256:db04c7114b656f896e206ba3873fe8d3a7adf2daa44907037f0274f4ba653fb9
ContainerName=forgejo
HostName=forgejo

Environment=USER_UID=1000
Environment=USER_GID=1000
Environment=TZ=Europe/Oslo

# Server
Environment=FORGEJO__server__DOMAIN=git.patriark.org
Environment=FORGEJO__server__ROOT_URL=https://git.patriark.org/
Environment=FORGEJO__server__HTTP_PORT=3000
# git-over-SSH: built-in server, published to loopback only (127.0.0.1:2222 below) so it's
# reachable from on-htpc SSH sessions (agent-forwarded keys) with NO LAN/internet exposure.
# Reverses the earlier "git-over-SSH deferred" decision; HTTPS path is unaffected.
Environment=FORGEJO__server__DISABLE_SSH=false
Environment=FORGEJO__server__START_SSH_SERVER=true
Environment=FORGEJO__server__SSH_LISTEN_PORT=2222
Environment=FORGEJO__server__SSH_PORT=2222
Environment=FORGEJO__server__SSH_DOMAIN=localhost

# Database (forgejo-db reached by container name on systemd-forgejo)
Environment=FORGEJO__database__DB_TYPE=postgres
Environment=FORGEJO__database__HOST=forgejo-db:5432
Environment=FORGEJO__database__NAME=forgejo
Environment=FORGEJO__database__USER=forgejo
Secret=forgejo_db_password,type=env,target=FORGEJO__database__PASSWD

# Closed forge (HTTPS-only, native auth; installer locked)
Environment=FORGEJO__security__INSTALL_LOCK=true
Environment=FORGEJO__service__DISABLE_REGISTRATION=true
Environment=FORGEJO__service__REQUIRE_SIGNIN_VIEW=true
Environment=FORGEJO__session__COOKIE_SECURE=true

# Storage (repos/LFS/data on container pool; not a DB, no NOCOW)
Volume=/mnt/btrfs-pool/subvol7-containers/forgejo:/data:Z

# Networks (reverse_proxy FIRST = default route for internet; static IPs per ADR-018)
Network=systemd-reverse_proxy:ip=10.89.2.89
Network=systemd-forgejo:ip=10.89.8.89

# git-over-SSH published to LOOPBACK ONLY — htpc-local; reached via agent-forwarded SSH sessions
PublishPort=127.0.0.1:2222:2222

DNS=192.168.1.69

# Health Check
HealthCmd=wget --quiet --tries=1 --spider http://localhost:3000/api/healthz || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryHigh=1G
MemoryMax=1536M

[Install]
WantedBy=default.target
