# quadlets/forgejo.container
[Unit]
Description=Forgejo Git Forge
After=network-online.target forgejo-db.service reverse_proxy-network.service forgejo-network.service
Wants=network-online.target
Requires=forgejo-db.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 15, resolved 2026-06-23 (index digest, multi-arch)
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=codeberg.org/forgejo/forgejo@sha256:55bb42bec9abef5223744804f164e37d37b20df7e8b8b4807ba213ad4f071d6d
ContainerName=forgejo
NoNewPrivileges=true
HostName=forgejo

Environment=USER_UID=1000
Environment=USER_GID=1000
Environment=TZ=Europe/Oslo

# Server
Environment=FORGEJO__server__DOMAIN=git.patriark.org
Environment=FORGEJO__server__ROOT_URL=https://git.patriark.org/
Environment=FORGEJO__server__HTTP_PORT=3000
# git-over-SSH: built-in server, published to loopback only (127.0.0.1:2222 below) so it's
# reachable from on-htpc SSH sessions (agent-forwarded keys) with NO LAN/internet exposure.
# Reverses the earlier "git-over-SSH deferred" decision; HTTPS path is unaffected.
Environment=FORGEJO__server__DISABLE_SSH=false
Environment=FORGEJO__server__START_SSH_SERVER=true
Environment=FORGEJO__server__SSH_LISTEN_PORT=2222
Environment=FORGEJO__server__SSH_PORT=2222
Environment=FORGEJO__server__SSH_DOMAIN=localhost

# Database (forgejo-db reached by container name on systemd-forgejo)
Environment=FORGEJO__database__DB_TYPE=postgres
Environment=FORGEJO__database__HOST=forgejo-db:5432
Environment=FORGEJO__database__NAME=forgejo
Environment=FORGEJO__database__USER=forgejo
Secret=forgejo_db_password,type=env,target=FORGEJO__database__PASSWD

# Closed forge (HTTPS-only, native auth; installer locked)
Environment=FORGEJO__security__INSTALL_LOCK=true
Environment=FORGEJO__service__DISABLE_REGISTRATION=true
Environment=FORGEJO__service__REQUIRE_SIGNIN_VIEW=true
Environment=FORGEJO__session__COOKIE_SECURE=true

# Instance commit signing (ADR-034): the FORGE signs its own server-side merge/
# squash commits so `main` stays fully verifiable. Authors sign their commits with
# their FIDO2 YubiKeys (ADR-033); a server can't tap a key, so the instance uses a
# passphraseless ed25519 SOFTWARE key. This closes the gap where Forgejo-created
# merge/squash commits (committer <name>@noreply.git.patriark.org) landed UNSIGNED.
#   - Section [repository.signing]: env escapes the section's dot as _0X2E_.
#   - SIGNING_KEY points at the PUBLIC key; ssh-keygen derives the private key by
#     stripping ".pub" (both files mounted below as podman secrets, same dir).
#   - SIGNING_EMAIL is the signer PRINCIPAL added to allowed_signers on author hosts.
#   - SIGNING_NAME contains a space -> the whole assignment MUST be double-quoted
#     (systemd Environment= splits unquoted values on spaces).
Environment=FORGEJO__repository_0X2E_signing__FORMAT=ssh
Environment=FORGEJO__repository_0X2E_signing__SIGNING_KEY=/run/secrets/forgejo_signing_key.pub
Environment="FORGEJO__repository_0X2E_signing__SIGNING_NAME=patriark.org Forgejo"
Environment=FORGEJO__repository_0X2E_signing__SIGNING_EMAIL=forgejo@git.patriark.org
Environment=FORGEJO__repository_0X2E_signing__MERGES=always

# Signing keypair as podman secrets, mounted owned by the in-container git user
# (uid/gid 1000 — gitea runs as uid 1000; podman maps these through the userns).
# Private key 0600, public 0644. Create the secrets BEFORE (re)starting — see ADR-034:
#   podman secret create forgejo_signing_key     ~/containers/secrets/forgejo_signing_key
#   podman secret create forgejo_signing_key_pub ~/containers/secrets/forgejo_signing_key.pub
Secret=forgejo_signing_key,type=mount,target=/run/secrets/forgejo_signing_key,uid=1000,gid=1000,mode=0600
Secret=forgejo_signing_key_pub,type=mount,target=/run/secrets/forgejo_signing_key.pub,uid=1000,gid=1000,mode=0644

# Storage (repos/LFS/data on container pool; not a DB, no NOCOW)
Volume=/mnt/btrfs-pool/subvol7-containers/forgejo:/data:Z

# Networks (reverse_proxy FIRST = default route for internet; static IPs per ADR-018)
Network=systemd-reverse_proxy:ip=10.89.2.89
Network=systemd-forgejo:ip=10.89.8.89

# git-over-SSH published to LOOPBACK ONLY — htpc-local; reached via agent-forwarded SSH sessions
PublishPort=127.0.0.1:2222:2222

DNS=192.168.1.69

# Health Check
HealthCmd=wget --quiet --tries=1 --spider http://localhost:3000/api/healthz || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
# Tier B (default 100) — ADR-035 boot tier: application class, deliberate default weight
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryHigh=1G
MemoryMax=1536M
CPUQuota=200%

[Install]
WantedBy=default.target
