# quadlets/gathio-db.container
[Unit]
Description=MongoDB Database for Gathio
After=network-online.target gathio-network.service
Wants=gathio-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 7, resolved 2026-05-23 (index digest, multi-arch)
Image=docker.io/library/mongo@sha256:32979a1189dfdc44da3f5ed40d910495f5ad8f6f7f77556646f890a30b2d3f56
ContainerName=gathio-db
HostName=gathio-db

# Database configuration (using podman secrets - Pattern 2: type=env)
Secret=gathio_mongodb_password,type=env,target=MONGO_INITDB_ROOT_PASSWORD
Environment=MONGO_INITDB_ROOT_USERNAME=gathio

# Storage (will set NOCOW attribute)
Volume=/mnt/btrfs-pool/subvol7-containers/gathio-db:/data/db:Z

# Network (dedicated Gathio network)
Network=systemd-gathio

# Health Check
HealthCmd=mongosh --quiet --eval "db.adminCommand('ping')" || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits (Best Practice: MemoryHigh + MemoryMax in [Service] section)
# MemoryHigh: soft limit - throttles when exceeded
# MemoryMax: hard limit - OOM kill
MemoryHigh=384M
MemoryMax=512M

[Install]
WantedBy=default.target
