# quadlets/gathio.container
[Unit]
Description=Gathio Event Management Platform
After=network-online.target gathio-db.service reverse_proxy-network.service monitoring-network.service
Wants=network-online.target
Requires=gathio-db.service

[Container]
Image=ghcr.io/lowercasename/gathio:latest
ContainerName=gathio
HostName=gathio
AutoUpdate=registry

# MongoDB connection (password mounted as file at /run/secrets/mongodb_password)
# config.toml will read this file for the connection string
Secret=gathio_mongodb_password,type=mount,target=/run/secrets/mongodb_password,mode=0400

# Gathio configuration
Environment=DOMAIN=events.patriark.org
Environment=PORT=3000
Environment=IS_FEDERATED=false
Environment=SHOW_KOFI=false

# Timezone
Environment=TZ=Europe/Oslo

# Storage volumes
Volume=%h/containers/config/gathio:/app/config:Z
Volume=%h/containers/data/gathio/static:/app/static:Z
Volume=%h/containers/data/gathio/images:/app/public/events:Z

# Networks (reverse_proxy FIRST for default route, then gathio for DB access)
Network=systemd-reverse_proxy
Network=systemd-gathio
Network=systemd-monitoring

# DNS
DNS=192.168.1.69

# Health Check
HealthCmd=wget --quiet --tries=1 --spider http://localhost:3000/ || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Slice=container.slice
Restart=always
TimeoutStartSec=300

# Resource Limits (Best Practice: MemoryHigh + MemoryMax in [Service] section)
# MemoryHigh: soft limit (512M) - throttles when exceeded
# MemoryMax: hard limit (768M) - OOM kill protection
MemoryHigh=512M
MemoryMax=768M

[Install]
WantedBy=default.target
