# quadlets/gathio-db.container
[Unit]
Description=MongoDB Database for Gathio
After=network-online.target gathio-network.service
Wants=gathio-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 7, resolved 2026-06-23 (index digest, multi-arch)
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=docker.io/library/mongo@sha256:8ecb514b00bdcc0bde67ef4e6c330385377a9dc68e24ee94e28c07c891647348
ContainerName=gathio-db
NoNewPrivileges=true
HostName=gathio-db

# Database configuration (using podman secrets - Pattern 2: type=env)
Secret=gathio_mongodb_password,type=env,target=MONGO_INITDB_ROOT_PASSWORD
Environment=MONGO_INITDB_ROOT_USERNAME=gathio

# Storage (will set NOCOW attribute)
Volume=/mnt/btrfs-pool/subvol8-db/gathio-db:/data/db:Z

# Network (dedicated Gathio network)
Network=systemd-gathio

# Health Check
HealthCmd=mongosh --quiet --eval "db.adminCommand('ping')" || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
# Boot tier A: elevated priority during default.target assembly (keystone)
StartupCPUWeight=200
StartupIOWeight=200

# Resource Limits (Best Practice: MemoryHigh + MemoryMax in [Service] section)
# MemoryHigh: soft limit - throttles when exceeded
# MemoryMax: hard limit - OOM kill
MemoryHigh=384M
MemoryMax=512M

[Install]
WantedBy=default.target
