# quadlets/grafana.container
[Unit]
Description=Grafana Monitoring Dashboard
After=network-online.target monitoring-network.service
Wants=network-online.target
Requires=monitoring-network.service

[Container]
Image=docker.io/grafana/grafana:latest
ContainerName=grafana
HostName=grafana
# Static IPs assigned to ensure consistent DNS resolution
Network=systemd-reverse_proxy:ip=10.89.2.9
Network=systemd-monitoring:ip=10.89.4.12

# Grafana Configuration (using podman secrets)
Environment=GF_SERVER_ROOT_URL=https://grafana.patriark.org
Environment=GF_SERVER_DOMAIN=grafana.patriark.org
Environment=GF_SECURITY_ADMIN_USER=patriark
Secret=grafana_admin_password,type=env,target=GF_SECURITY_ADMIN_PASSWORD
Environment=GF_AUTH_DISABLE_LOGIN_FORM=false
Environment=GF_AUTH_PROXY_ENABLED=true
Environment=GF_AUTH_PROXY_HEADER_NAME=X-Forwarded-User
Environment=GF_AUTH_PROXY_HEADER_PROPERTY=username
Environment=GF_AUTH_PROXY_AUTO_SIGN_UP=true

# Storage
Volume=%h/containers/data/grafana:/var/lib/grafana:Z
Volume=%h/containers/config/grafana/provisioning:/etc/grafana/provisioning:ro,Z

DNS=192.168.1.69

# Health check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
AutoUpdate=registry

[Service]
Slice=container.slice
MemoryMax=1G
MemoryHigh=900M
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=default.target
