# quadlets/home-assistant.container
[Unit]
Description=Home Assistant - Home Automation Hub
After=network-online.target reverse_proxy-network.service home_automation-network.service monitoring-network.service
Wants=network-online.target
Requires=reverse_proxy-network.service home_automation-network.service monitoring-network.service

[Container]
Image=ghcr.io/home-assistant/home-assistant:stable
ContainerName=home-assistant
HostName=home-assistant
AutoUpdate=registry

# Networks (reverse_proxy FIRST for internet access)
# Static IPs assigned to ensure consistent routing (see ADR-018)
Network=systemd-reverse_proxy:ip=10.89.2.76
Network=systemd-home_automation:ip=10.89.6.76
Network=systemd-monitoring:ip=10.89.4.76

# Environment
Environment=TZ=Europe/Oslo

# DNS
DNS=192.168.1.69
DNSSearch=lokal

# Storage Volumes
Volume=%h/containers/config/home-assistant:/config:Z
Volume=%h/containers/data/home-assistant:/data:Z

# Ports (published for local access)
PublishPort=8123:8123

# Health Check
HealthCmd=curl -f http://localhost:8123/manifest.json || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=90s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryMax=2G
MemoryHigh=1800M
MemorySwapMax=512M

[Install]
WantedBy=default.target
