# quadlets/home-assistant.container
[Unit]
Description=Home Assistant - Home Automation Hub
After=network-online.target reverse_proxy-network.service home_automation-network.service
Wants=home_automation-network.service network-online.target reverse_proxy-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: stable, resolved 2026-05-23 (index digest, multi-arch)
Image=ghcr.io/home-assistant/home-assistant@sha256:d4fbec16196d5c8bedf32647f0ca7165f654d92a2e81f35a373508d3226cb867
ContainerName=home-assistant
NoNewPrivileges=true
HostName=home-assistant

# Networks (reverse_proxy FIRST for internet access)
# Static IPs assigned to ensure consistent routing (see ADR-018)
Network=systemd-reverse_proxy:ip=10.89.2.76
Network=systemd-home_automation:ip=10.89.6.76

# Environment
Environment=TZ=Europe/Oslo

# DNS
DNS=192.168.1.69
DNSSearch=lokal

# Storage Volumes
Volume=%h/containers/config/home-assistant:/config:Z
Volume=%h/containers/data/home-assistant:/data:Z

# Ports (bound to LAN interface only - see GH #142)
PublishPort=192.168.1.70:8123:8123

# Health Check
HealthCmd=curl -f http://localhost:8123/manifest.json || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=90s

[Service]
# Tier B (default 100) — ADR-035 boot tier: application class, deliberate default weight
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryMax=2G
MemoryHigh=1800M
MemorySwapMax=512M
CPUQuota=400%

[Install]
WantedBy=default.target
