# quadlets/immich-server.container
[Unit]
Description=Immich Server
After=network-online.target photos-network.service postgresql-immich.service redis-immich.service
Wants=network-online.target
Requires=photos-network.service postgresql-immich.service redis-immich.service

[Container]
Image=ghcr.io/immich-app/immich-server:v2.5.5
ContainerName=immich-server
AutoUpdate=registry

# Run as non-root user (security: ADR-001 Rootless Containers)
# NOTE: User=1000:1000 removed due to Immich folder integrity check incompatibility
# Security provided by rootless Podman (UID 0 inside → UID 1000 outside) + Traefik layers
# TODO: Revisit after consulting Immich community or upstream fix

# Networks (multiple - order matters for default route)
# Static IPs assigned to ensure consistent DNS resolution
Network=systemd-reverse_proxy:ip=10.89.2.12
Network=systemd-photos:ip=10.89.5.5
Network=systemd-monitoring:ip=10.89.4.22

# Environment - Database
Environment=DB_HOSTNAME=postgresql-immich
Environment=DB_USERNAME=immich
Environment=DB_DATABASE_NAME=immich
Secret=postgres-password,type=env,target=DB_PASSWORD

# Environment - Redis
Environment=REDIS_HOSTNAME=redis-immich
Environment=REDIS_PORT=6379

# Environment - Machine Learning
Environment=IMMICH_MACHINE_LEARNING_URL=http://immich-ml:3003

# Environment - Security
Secret=immich-jwt-secret,type=env,target=JWT_SECRET

# Environment - Upload
Environment=UPLOAD_LOCATION=/usr/src/app/upload

# Storage - Photo library on BTRFS pool
Volume=/mnt/btrfs-pool/subvol3-opptak/immich:/usr/src/app/upload:Z

# External library - read-write so Immich can manage (delete/move) assets directly
Volume=/mnt/btrfs-pool/subvol3-opptak/immich/library:/mnt/media:Z

# Hardware Acceleration - GPU access for VAAPI transcoding
AddDevice=/dev/dri:/dev/dri

# Port exposure: Use Traefik exclusively (no direct port binding)
# Removed PublishPort - access via https://photos.patriark.org only
# PublishPort=2283:2283

# Health check
HealthCmd=curl -f http://localhost:2283/api/server/ping || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=300s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=900

# Resource Limits
MemoryMax=4G
MemoryHigh=3.6G
MemorySwapMax=1G

[Install]
WantedBy=default.target
