# quadlets/immich-ml.container
[Unit]
Description=Immich Machine Learning
After=network-online.target photos-network.service
Wants=network-online.target
Requires=photos-network.service

[Container]
Image=ghcr.io/immich-app/immich-machine-learning:v2.5.5
ContainerName=immich-ml
AutoUpdate=registry

# Network (isolated - only on photos network)
Network=systemd-photos

# Storage - ML model cache on BTRFS pool
Volume=/mnt/btrfs-pool/subvol7-containers/immich-ml-cache:/cache:Z

# Environment
Environment=MACHINE_LEARNING_CACHE_FOLDER=/cache

# Resource limits (podman runtime)
Memory=4G

# Health check
# Using python3 for IPv4/IPv6 compatibility (service listens on [::]:3003)
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:3003/ping', timeout=5)" || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=600s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=1800

# Resource Limits (systemd cgroup)
MemoryMax=4G
MemoryHigh=3.6G
MemorySwapMax=512M

[Install]
WantedBy=default.target
