# quadlets/loki.container
[Unit]
Description=Loki Log Aggregation System
After=network-online.target monitoring-network.service
Wants=network-online.target
Requires=monitoring-network.service

[Container]
Image=docker.io/grafana/loki:3.2.1
ContainerName=loki
HostName=loki
Network=systemd-reverse_proxy
Network=systemd-monitoring

# Loki configuration and data
Volume=%h/containers/config/loki/loki-config.yml:/etc/loki/local-config.yaml:ro,Z
Volume=/mnt/btrfs-pool/subvol7-containers/loki:/loki:Z

# Loki command line arguments
Exec=-config.file=/etc/loki/local-config.yaml

DNS=192.168.1.69

# Health check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:3100/ready || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
AutoUpdate=registry

[Service]
Restart=on-failure
TimeoutStartSec=120

# Resource Limits
MemoryMax=512M

[Install]
WantedBy=default.target
