# quadlets/loki.container
[Unit]
Description=Loki Log Aggregation System
After=network-online.target monitoring-network.service
Wants=network-online.target
Requires=monitoring-network.service

[Container]
Image=docker.io/grafana/loki:latest
ContainerName=loki
HostName=loki
# Static IPs assigned to ensure consistent DNS resolution
Network=systemd-reverse_proxy:ip=10.89.2.15
Network=systemd-monitoring:ip=10.89.4.23

# Loki configuration and data
Volume=%h/containers/config/loki/loki-config.yml:/etc/loki/local-config.yaml:ro,Z
Volume=/mnt/btrfs-pool/subvol7-containers/loki:/loki:Z

# Loki command line arguments
Exec=-config.file=/etc/loki/local-config.yaml

DNS=192.168.1.69

# No container healthcheck possible - Loki uses a scratch/distroless image with no shell or tools
# External monitoring via Prometheus up{job="loki"} and homelab-intel.sh (verifies via Promtail logs)
AutoUpdate=registry

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=120

# Resource Limits
MemoryMax=512M
MemoryHigh=460M

[Install]
WantedBy=default.target
