# quadlets/nextcloud-db.container
[Unit]
Description=Nextcloud MariaDB Database
After=network-online.target
Wants=network-online.target


[Container]
Image=docker.io/library/mariadb:11
ContainerName=nextcloud-db
AutoUpdate=registry

# MariaDB configuration (using podman secrets)
Secret=nextcloud_db_root_password,type=env,target=MYSQL_ROOT_PASSWORD
Environment=MYSQL_DATABASE=nextcloud
Environment=MYSQL_USER=nextcloud
Secret=nextcloud_db_password,type=env,target=MYSQL_PASSWORD
Environment=MYSQL_INNODB_FLUSH_METHOD=O_DIRECT
Environment=MYSQL_INNODB_FILE_PER_TABLE=1

# Storage (BTRFS with NOCOW)
Volume=/mnt/btrfs-pool/subvol7-containers/nextcloud-db/data:/var/lib/mysql:Z

# Networks
Network=systemd-nextcloud
Network=systemd-monitoring

# Health check
HealthCmd=healthcheck.sh --connect --innodb_initialized
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5

# Labels for Prometheus
Label=prometheus.scrape=true
Label=prometheus.port=3306

[Service]
Restart=on-failure
TimeoutStartSec=300

# Resource limits (systemd memory control)
MemoryMax=1G
MemoryHigh=900M

[Install]
WantedBy=default.target
