# quadlets/nextcloud-db.container
[Unit]
Description=Nextcloud MariaDB Database
After=network-online.target
Wants=network-online.target


[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 11, resolved 2026-07-13 (index digest, multi-arch)
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=docker.io/library/mariadb@sha256:efb4959ef2c835cd735dbc388eb9ad6aab0c78dd64febcd51bc17481111890c4
ContainerName=nextcloud-db
NoNewPrivileges=true

# MariaDB configuration (using podman secrets)
Secret=nextcloud_db_root_password,type=env,target=MYSQL_ROOT_PASSWORD
Environment=MYSQL_DATABASE=nextcloud
Environment=MYSQL_USER=nextcloud
Secret=nextcloud_db_password,type=env,target=MYSQL_PASSWORD
Environment=MYSQL_INNODB_FLUSH_METHOD=O_DIRECT
Environment=MYSQL_INNODB_FILE_PER_TABLE=1

# Storage (BTRFS with NOCOW)
Volume=/mnt/btrfs-pool/subvol8-db/nextcloud-db:/var/lib/mysql:Z

# Networks
Network=systemd-nextcloud

# Health check
HealthCmd=healthcheck.sh --connect --innodb_initialized
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5

# Labels for Prometheus
Label=prometheus.scrape=true
Label=prometheus.port=3306

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
OOMPolicy=kill
# Boot tier A: elevated priority during default.target assembly (keystone)
StartupCPUWeight=200
StartupIOWeight=200

# Resource limits (systemd memory control)
MemoryMax=1G
MemoryHigh=900M

[Install]
WantedBy=default.target
