# quadlets/node_exporter.container
[Unit]
Description=Node Exporter - System Metrics Collector
After=network-online.target monitoring-network.service traefik.service
Wants=monitoring-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-23 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=quay.io/prometheus/node-exporter@sha256:0f422f62c15f154af8d8572b23d623aebfb10cec73a5c654d18f911f3f9df241
ContainerName=node_exporter
NoNewPrivileges=true
HostName=node_exporter
Network=systemd-monitoring

# rshared propagation on / so removable-drive unmounts propagate in — see ADR-023.
Volume=/:/host:ro,rshared
Volume=/sys:/host/sys:ro,rslave
Volume=/proc:/host/proc:ro,rslave

# Mount textfile collector directory for custom metrics (backup, restore test)
Volume=%h/containers/data/backup-metrics:/textfile-metrics:ro,Z

# Tell node_exporter where to find host filesystem
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Node exporter command line arguments
# mount-points-exclude extends the upstream default with run/media/ for removable backup drives (ADR-023)
Exec=--path.rootfs=/host --collector.textfile.directory=/textfile-metrics --collector.filesystem.mount-points-exclude=^/(dev|proc|run/(credentials/.+|media/.+)|sys|var/lib/(docker|containers)/.+|sysroot)($|/)

DNS=192.168.1.69

# Health check - use proper GET request (--spider uses HEAD which /metrics doesn't support)
HealthCmd=wget --no-verbose --tries=1 -O /dev/null http://localhost:9100/metrics || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=60
# Boot tier C: yield during fan-out; After=traefik.service defers past the storm
StartupCPUWeight=50
StartupIOWeight=50

# Resource Limits
MemoryMax=128M
MemoryHigh=115M

[Install]
WantedBy=default.target
