# quadlets/node_exporter.container
[Unit]
Description=Node Exporter - System Metrics Collector
After=network-online.target monitoring-network.service
Wants=network-online.target
Requires=monitoring-network.service

[Container]
Image=quay.io/prometheus/node-exporter:latest
ContainerName=node_exporter
HostName=node_exporter
Network=systemd-monitoring

# Node exporter needs access to host filesystem for accurate metrics
# Mount host paths as read-only
Volume=/:/host:ro,rslave
Volume=/sys:/host/sys:ro,rslave
Volume=/proc:/host/proc:ro,rslave

# Mount textfile collector directory for custom metrics (backup, restore test)
Volume=%h/containers/data/backup-metrics:/textfile-metrics:ro,Z

# Tell node_exporter where to find host filesystem
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Node exporter command line arguments
Exec=--path.rootfs=/host --collector.textfile.directory=/textfile-metrics

DNS=192.168.1.69

# Resource limits (podman runtime)
Memory=128M

# Health check - use proper GET request (--spider uses HEAD which /metrics doesn't support)
HealthCmd=wget --no-verbose --tries=1 -O /dev/null http://localhost:9100/metrics || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
AutoUpdate=registry

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=60

# Resource Limits
MemoryMax=128M
MemoryHigh=115M

[Install]
WantedBy=default.target
