# quadlets/ocis.container
# ~/.config/containers/systemd/ocis.container
# ownCloud Infinite Scale (OCIS) - Modern Cloud Platform
# Created: 2025-11-12
# Documentation: https://doc.owncloud.com/ocis/

# ═══════════════════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════════════════
[Unit]
Description=ownCloud Infinite Scale (OCIS)
Documentation=https://doc.owncloud.com/ocis/
After=network-online.target traefik.service
Wants=network-online.target
Requires=traefik.service

# ═══════════════════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════════════════
[Container]
Image=docker.io/owncloud/ocis:7.1.3
ContainerName=ocis
HostName=ocis

# Networks
Network=reverse_proxy.network
Network=monitoring.network

# Volumes - OCIS Data & Config
# Note: subvol1-docs and subvol2-pics NOT mounted
# Decision: Use OCIS native storage for better performance and features
# Original data remains available for import if needed
Volume=/mnt/btrfs-pool/subvol7-containers/ocis/data:/var/lib/ocis:Z
Volume=%h/containers/config/ocis:/etc/ocis:Z

# Environment - Core Configuration
Environment=OCIS_URL=https://cloud.patriark.org
Environment=OCIS_INSECURE=false
Environment=OCIS_LOG_LEVEL=info
Environment=OCIS_LOG_PRETTY=false

# Environment - Proxy Configuration (behind Traefik)
Environment=PROXY_HTTP_ADDR=0.0.0.0:9200
Environment=PROXY_TLS=false

# Environment - Admin Account (change on first login!)
Environment=OCIS_ADMIN_USER_ID=admin
Environment=OCIS_ADMIN_USER_PASSWORD=CHANGE_ME_ON_FIRST_LOGIN

# Environment - Storage
Environment=STORAGE_USERS_OCIS_ROOT=/var/lib/ocis

# Environment - Metrics/Debug (for Prometheus scraping)
Environment=PROXY_DEBUG_ADDR=0.0.0.0:9205

# Secrets (Podman secrets mapped to environment variables)
Secret=ocis_jwt_secret,type=env,target=OCIS_JWT_SECRET
Secret=ocis_transfer_secret,type=env,target=OCIS_TRANSFER_SECRET
Secret=ocis_machine_auth_api_key,type=env,target=OCIS_MACHINE_AUTH_API_KEY

# Health check
HealthCmd=curl -f http://localhost:9200/status.php || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# ═══════════════════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════════════════
[Service]
Restart=on-failure
TimeoutStartSec=120
RestartSec=10s

# Resource Limits
MemoryMax=2G
MemoryHigh=1.5G

# ═══════════════════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════════════════
[Install]
WantedBy=default.target
