# quadlets/pihole-exporter.container
[Unit]
Description=Pi-hole Exporter (Prometheus) — query/block metrics from the LAN resolver
Documentation=https://github.com/eko/pihole-exporter
After=network-online.target restricted-egress-network.service monitoring-network.service traefik.service
Wants=network-online.target restricted-egress-network.service monitoring-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-27.
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only.
# Pi-hole v6: uses app-password auth, supplied via Podman secret PIHOLE_PASSWORD.
Image=docker.io/ekofr/pihole-exporter@sha256:a890cc731a39da719ba19cfd717427d93366cd36307f069d495f37a1373ccffc
ContainerName=pihole-exporter
HostName=pihole-exporter

# restricted-egress first → default route to LAN (the Pi at 192.168.1.69),
# internet egress dropped by the netns nft filter (ADR-045). Single-homing
# invariant: this must stay the ONLY non-internal network on this container.
# monitoring → Prometheus scrapes :9617 (pinned IP; name resolves ambiguously
# across the two shared networks post-swap).
Network=systemd-restricted-egress
Network=systemd-monitoring:ip=10.89.4.96

Environment=PIHOLE_HOSTNAME=192.168.1.69
Environment=PIHOLE_PORT=80
Environment=PIHOLE_PROTOCOL=http
Environment=PORT=9617
Environment=INTERVAL=30s
Environment=TZ=Europe/Oslo

# Pi-hole v6 app password — generate in Pi-hole UI (Settings → Web Interface / API → App passwords).
# Create secret:  printf '<app-password>' | podman secret create pihole_api_token -
Secret=pihole_api_token,type=env,target=PIHOLE_PASSWORD

# NO HealthCmd: this image (ekofr/pihole-exporter) is distroless — it has no
# wget/curl/sh, so the previous `HealthCmd=wget ...` could NEVER run ("executable
# file wget not found"). The container therefore reported a PERMANENT false
# "unhealthy", and because the process never exits, Restart=on-failure could not
# auto-heal a wedged-but-running exporter. This is exactly how the 2026-06-04
# latency regression went unrecovered for 3 days. Liveness/self-heal is now
# enforced externally by pihole-exporter-watchdog.timer (host-side probe of
# /metrics with restart-on-stall), which does not depend on in-container tooling.

User=1000:1000
NoNewPrivileges=true

[Service]
Slice=container.slice
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s
# Boot tier C: yield during fan-out; After=traefik.service defers past the storm
StartupCPUWeight=50
StartupIOWeight=50

MemoryMax=64M
MemoryHigh=56M
CPUQuota=25%

[Install]
WantedBy=default.target
