# quadlets/postgres-exporter.container
[Unit]
Description=PostgreSQL Exporter for postgresql-immich (Prometheus)
After=network-online.target photos-network.service monitoring-network.service postgresql-immich.service
Wants=network-online.target photos-network.service monitoring-network.service
Requires=postgresql-immich.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-07-13 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = adopt new digest (pin-container-image.sh --adopt), bake, restart.
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=quay.io/prometheuscommunity/postgres-exporter@sha256:ac5ec343104fae0e2d84a27bb8d69b38430a11910c5382cad85d478d2bab713e
ContainerName=postgres-exporter
NoNewPrivileges=true
HostName=postgres-exporter

# Reach postgresql-immich on photos; expose metrics on monitoring.
Network=systemd-monitoring:ip=10.89.4.94
Network=systemd-photos:ip=10.89.5.94

# Credentials: dedicated prometheus_exporter role with pg_monitor (created out-of-band, not in quadlet).
Secret=postgres-exporter-password,type=env,target=DATA_SOURCE_PASS
Environment=DATA_SOURCE_URI=postgresql-immich:5432/immich?sslmode=disable
Environment=DATA_SOURCE_USER=prometheus_exporter

# L-018: custom query exposing immich_active_assets (mass-deletion canary).
# Grant requirements + rationale documented in the queries file itself.
Volume=%h/containers/config/postgres-exporter/queries.yml:/etc/postgres-exporter/queries.yml:ro,Z
Environment=PG_EXPORTER_EXTEND_QUERY_PATH=/etc/postgres-exporter/queries.yml

HealthCmd=wget --no-verbose --tries=1 -O /dev/null http://localhost:9187/metrics || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=60
# Boot tier C: yield during fan-out (already gated behind postgresql-immich.service)
StartupCPUWeight=50
StartupIOWeight=50

MemoryMax=128M
MemoryHigh=115M

[Install]
WantedBy=default.target
