# quadlets/prometheus.container
[Unit]
Description=Prometheus Metrics Database
After=network-online.target monitoring-network.service node_exporter.service traefik.service
Wants=monitoring-network.service network-online.target node_exporter.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-06-10 (index digest, multi-arch)
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=quay.io/prometheus/prometheus@sha256:69f5241418838263316593f7274a304b095c40bcf22e57272865da91bd60a8ac
ContainerName=prometheus
NoNewPrivileges=true
HostName=prometheus
# Static IPs assigned to ensure consistent DNS resolution
Network=systemd-reverse_proxy:ip=10.89.2.79
Network=systemd-monitoring:ip=10.89.4.79

# Prometheus configuration and data
# Directory mount (GH#276): single-file binds pin the host inode at start, so
# atomic-replace writes (editor save, sed -i, git checkout) silently orphan the
# bind and /-/reload applies nothing. Bit 4 times; directory mounts don't.
Volume=%h/containers/config/prometheus:/etc/prometheus:ro,Z
Volume=/mnt/btrfs-pool/subvol8-db/prometheus:/prometheus:Z

# Podman secrets for authentication
Secret=ha-prometheus-token,type=mount,target=/run/secrets/ha-prometheus-token

# Prometheus command line arguments
# --web.enable-admin-api exposes the TSDB snapshot endpoint used by db-dump.sh
# (ADR-029 Tier-2 backup). It also exposes destructive endpoints
# (/api/v1/admin/tsdb/delete_series, clean_tombstones) to the unauthenticated
# systemd-monitoring network. Accepted for this homelab: the external path is
# Authelia-fronted (routers.yml), only the snapshot endpoint is used, and the
# dump job pre-cleans orphaned snapshots. See ADR-029.
#
# --enable-feature=memory-snapshot-on-shutdown: on a clean SIGTERM (the path
# graceful-shutdown.sh / prepare-for-reboot.sh always take) Prometheus writes a
# head-chunk snapshot and loads it on next start INSTEAD of replaying the WAL
# (~16s of HDD reads on subvol8-db, among the heaviest fleet boot I/O — the very
# reason this service sits in boot tier C, see ADR-035). Falls back to WAL replay
# on an unclean stop, so it never hurts. GH#308.
Exec=--config.file=/etc/prometheus/prometheus.yml \
     --storage.tsdb.path=/prometheus \
     --storage.tsdb.retention.time=15d \
     --enable-feature=memory-snapshot-on-shutdown \
     --web.console.libraries=/usr/share/prometheus/console_libraries \
     --web.console.templates=/usr/share/prometheus/consoles \
     --web.enable-lifecycle \
     --web.enable-admin-api

DNS=192.168.1.69

# Health check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:9090/-/healthy || exit 1
HealthInterval=10s
HealthTimeout=3s
HealthRetries=3

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=120
OOMPolicy=kill
# Boot tier C: yield during fan-out; After=traefik.service defers past the storm.
# Added 2026-06-09: TSDB mmap+WAL replay (~20s of reads on subvol8-db/HDD) is among
# the heaviest fleet boot I/O. Previously implicit tier B (default 100) — reclassified
# after a recurring boot storm. See ADR-035.
StartupCPUWeight=50
StartupIOWeight=50

# Resource Limits
MemoryMax=1G
MemoryHigh=900M
CPUQuota=400%

[Install]
WantedBy=default.target
