# quadlets/qbittorrent.container
[Unit]
Description=qBittorrent - BitTorrent Client (Headless)
After=network-online.target reverse_proxy-network.service
Wants=network-online.target reverse_proxy-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-23 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = resolve new digest, bake, edit, restart.
Image=docker.io/linuxserver/qbittorrent@sha256:f76c4363cce0834df43b097682928d217f9f0e0fa1d13781c0934848d543e4ab
ContainerName=qbittorrent
HostName=qbittorrent
Network=systemd-reverse_proxy:ip=10.89.2.85
# No PublishPort - WebUI (port 8085, set via WEBUI_PORT below) accessible only
# via Traefik reverse proxy. No BT port published - Session\Port=30069 is set
# in config/qbittorrent/qBittorrent/qBittorrent.conf but qBt operates in
# passive/NAT mode (outbound connections only), so no UDM forward is required.
# NOTE: `podman port qbittorrent` reports 6881/tcp, 6881/udp, 8080/tcp from the
# upstream image's EXPOSE directives - those are NOT the actual ports in use.

# linuxserver.io PUID/PGID
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Europe/Oslo
Environment=WEBUI_PORT=8085

# Config
Volume=%h/containers/config/qbittorrent:/config:Z

# Media directories - mounted at same paths as host to match fastresume save paths
Volume=/mnt/btrfs-pool/subvol4-multimedia:/mnt/btrfs-pool/subvol4-multimedia:z
Volume=/mnt/btrfs-pool/subvol6-tmp/Downloads:/mnt/btrfs-pool/subvol6-tmp/Downloads:z

DNS=192.168.1.69

# Health check
HealthCmd=curl -f http://localhost:8085 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=120

# Resource Limits
MemoryMax=1G
MemoryHigh=900M
CPUQuota=200%

[Install]
WantedBy=default.target
