# quadlets/redis-authelia.container
[Unit]
Description=Redis - Authelia Session Storage
After=network-online.target auth_services-network.service
Wants=auth_services-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: 7-alpine, resolved 2026-05-23 (index digest, multi-arch)
Image=docker.io/library/redis@sha256:6ab0b6e7381779332f97b8ca76193e45b0756f38d4c0dcda72dbb3c32061ab99
ContainerName=redis-authelia
NoNewPrivileges=true

# Network
Network=systemd-auth_services

# Storage for session persistence
Volume=%h/containers/data/redis-authelia:/data:Z

# Redis configuration for session storage
# volatile-lru: only evict keys with TTLs (sessions). Regulation counters
# (no TTL) are preserved across pressure — brute-force protection survives.
Exec=redis-server --appendonly yes --appendfsync everysec --maxmemory 200mb --maxmemory-policy volatile-lru

# Health check
HealthCmd=redis-cli ping
HealthInterval=10s
HealthTimeout=3s
HealthRetries=5

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
OOMPolicy=kill
# Boot tier A: elevated priority during default.target assembly (keystone)
StartupCPUWeight=200
StartupIOWeight=200

# Resource Limits
MemoryMax=256M
MemoryHigh=230M

[Install]
WantedBy=default.target
