# quadlets/redis-immich.container
[Unit]
Description=Redis for Immich
After=network-online.target photos-network.service
Wants=network-online.target photos-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-06-10 (index digest, multi-arch)
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=docker.io/valkey/valkey@sha256:4963247afc4cd33c7d3b2d2816b9f7f8eeebab148d29056c2ca4d7cbc966f2d9
ContainerName=redis-immich
NoNewPrivileges=true

# Network
Network=systemd-photos

# Storage (BTRFS pool with NOCOW for persistence)
Volume=/mnt/btrfs-pool/subvol7-containers/redis-immich:/data:Z

# Health check
HealthCmd=valkey-cli ping
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300
# Boot tier A: elevated priority during default.target assembly (keystone)
StartupCPUWeight=200
StartupIOWeight=200

# Resource Limits
MemoryMax=512M
MemoryHigh=460M

[Install]
WantedBy=default.target
