# quadlets/redis-authelia.container
[Unit]
Description=Redis - Authelia Session Storage
After=network-online.target auth_services-network.service
Wants=network-online.target
Requires=auth_services-network.service

[Container]
Image=docker.io/library/redis:7-alpine
ContainerName=redis-authelia
AutoUpdate=registry

# Network
Network=systemd-auth_services

# Storage for session persistence
Volume=%h/containers/data/redis-authelia:/data:Z

# Redis configuration for session storage
Exec=redis-server --appendonly yes --appendfsync everysec --maxmemory 128mb --maxmemory-policy allkeys-lru

# Health check
HealthCmd=redis-cli ping
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryMax=256M
MemoryHigh=230M

[Install]
WantedBy=default.target
