# quadlets/redis-immich.container
[Unit]
Description=Redis for Immich
After=network-online.target photos-network.service
Wants=network-online.target photos-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-05-23 (index digest, multi-arch)
Image=docker.io/valkey/valkey@sha256:8436e10bc65c94886a91d4415b6a6dfa9cb5a306fb3b996e5bb67cd2b4854193
ContainerName=redis-immich

# Network
Network=systemd-photos

# Storage (BTRFS pool with NOCOW for persistence)
Volume=/mnt/btrfs-pool/subvol7-containers/redis-immich:/data:Z

# Health check
HealthCmd=valkey-cli ping
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Slice=container.slice
Restart=on-failure
TimeoutStartSec=300

# Resource Limits
MemoryMax=512M
MemoryHigh=460M

[Install]
WantedBy=default.target
