# quadlets/tinyauth.container
[Unit]
Description=TinyAuth - Lightweight Forward Authentication
After=traefik.service reverse_proxy-network.service auth_services-network.service
Wants=reverse_proxy-network.service auth_services-network.service
Requires=reverse_proxy-network.service auth_services-network.service

[Container]
Image=ghcr.io/steveiliop56/tinyauth:v4
ContainerName=tinyauth
HostName=tinyauth

# Networks (multi-network: serves login page + handles auth checks)
Network=systemd-reverse_proxy
Network=systemd-auth_services

# ⚠️  REPLACE THESE PLACEHOLDERS BEFORE DEPLOYMENT ⚠️
Environment=APP_URL=https://auth.patriark.org
Environment=SECRET=#this-is-a-placeholder
Environment=USERS=#this-is-another-placeholder

# Storage
Volume=%h/containers/data/tinyauth:/app/data:Z

# Health check - validates web server is responding (checks login page)
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:3000/ || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=on-failure
TimeoutStartSec=60

# Resource Limits
MemoryMax=256M

[Install]
WantedBy=default.target
