# quadlets/unifi-syslog.container
[Unit]
Description=UniFi UDM Pro Syslog Receiver (syslog-ng)
After=network-online.target syslog-network.service traefik.service
Wants=syslog-network.service network-online.target

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-06-10 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = adopt new digest (pin-container-image.sh --adopt), bake, restart.
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=docker.io/linuxserver/syslog-ng@sha256:b77c32d93b9e9b84f4fdd9261d7dd1db3b8a6f90391459b1ca51c82af9828bf8
ContainerName=unifi-syslog
HostName=unifi-syslog
Network=systemd-syslog:ip=10.89.7.69

# Bind to LAN interface only - NOT 0.0.0.0. Defense in depth alongside firewalld rule.
PublishPort=192.168.1.70:1514:1514/udp

# Rootless Podman maps container-uid-0 to host-uid-1000 (patriark), which owns the bind
# mount. Running syslog-ng as container-root keeps write semantics simple and avoids the
# LSIO abc-user subuid shuffle that requires either a :U volume flag (unreliable with :Z)
# or host-side subuid chown. For a write-only forensic-log receiver this is the cleaner
# pattern.
Environment=PUID=0
Environment=PGID=0
Environment=TZ=Europe/Oslo

Volume=%h/containers/config/unifi-syslog/syslog-ng.conf:/config/syslog-ng.conf:ro,Z
# Tier 1 (subvol7): snapshotted + sent offsite by Urd. Moved off subvol8-db
# (ADR-029): append-only forensic logs belong in the snapshot tier, and the
# COW cost is negligible for append-only writes.
Volume=/mnt/btrfs-pool/subvol7-containers/unifi-syslog:/var/log/unifi-syslog:Z

DNS=192.168.1.69

HealthCmd=pgrep -x syslog-ng
HealthInterval=60s
HealthRetries=3
HealthStartPeriod=30s
HealthTimeout=5s

[Service]
Slice=container.slice
Restart=on-failure
# Boot-storm headroom: see blackbox-exporter.container for the cold-reboot incident
# 2026-05-27 where 60s timeout fired mid-startup → ExecStopPost SIGABRT.
TimeoutStartSec=180
TimeoutStopSec=60
# Boot tier C: yield during fan-out; After=traefik.service defers past the storm
StartupCPUWeight=50
StartupIOWeight=50

# Resource Limits
MemoryMax=128M
MemoryHigh=110M

[Install]
WantedBy=default.target
