# quadlets/unpoller.container
[Unit]
Description=Unpoller - UniFi Metrics Exporter for Prometheus
Documentation=https://unpoller.com/docs/
After=network-online.target restricted-egress-network.service prometheus.service traefik.service
Wants=network-online.target restricted-egress-network.service

[Container]
# ADR-030 P2: digest-pinned (integrity). tag: latest, resolved 2026-07-13 (index digest, multi-arch)
# ADR-030 P1/P4: AutoUpdate removed — deliberate updates only. Bump = adopt new digest (pin-container-image.sh --adopt), bake, restart.
# ADR-030 P6: no publisher signature (tracked unsigned)
Image=ghcr.io/unpoller/unpoller@sha256:9dcccdc931a6830735f6978caf8cd67699b0dc33e37cf9ef4638611791c4df62
ContainerName=unpoller

# Configuration
Volume=%h/containers/config/unpoller/up.conf:/etc/unpoller/up.conf:Z,ro
Volume=%h/containers/data/unpoller:/data:Z

# Environment variables for UniFi controller credentials
# These should be set via podman secret or systemd EnvironmentFile
# For now, using placeholders - update before starting
Secret=unpoller_url,type=env,target=UP_UNIFI_DEFAULT_URL
Secret=unpoller_user,type=env,target=UP_UNIFI_DEFAULT_USER
Secret=unpoller_pass,type=env,target=UP_UNIFI_DEFAULT_PASS
Environment=TZ=Europe/Oslo

# Network configuration
# restricted-egress first for default route (UDM Pro polling at 192.168.1.1);
# internet egress dropped by the netns nft filter (ADR-045). Single-homing
# invariant: this must stay the ONLY non-internal network on this container.
# monitoring IP pinned — Prometheus scrapes 10.89.4.97:9130 (name resolves
# ambiguously across the two shared networks post-swap).
Network=systemd-restricted-egress
Network=systemd-monitoring:ip=10.89.4.97

# Prometheus metrics exposed on systemd-monitoring network (port 9130)
# No PublishPort needed - Prometheus can reach unpoller:9130 directly

# Health check disabled - distroless container lacks curl
# Monitoring via Prometheus scrape health instead
# HealthCmd=curl -f http://localhost:9130/metrics || exit 1
# HealthInterval=30s
# HealthRetries=3
# HealthStartPeriod=10s
# HealthTimeout=5s

# Security
User=1000:1000
NoNewPrivileges=true
SecurityLabelDisable=false

[Service]
Slice=container.slice

# Resource limits
MemoryMax=256M
MemoryHigh=230M
CPUQuota=50%

# Restart policy
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s
# Boot tier C: yield during fan-out; After=traefik.service defers past the storm
StartupCPUWeight=50
StartupIOWeight=50

[Install]
WantedBy=default.target
