# quadlets/unpoller.container
[Unit]
Description=Unpoller - UniFi Metrics Exporter for Prometheus
Documentation=https://unpoller.com/docs/
After=network-online.target prometheus.service
Wants=network-online.target

[Container]
Image=ghcr.io/unpoller/unpoller:latest
ContainerName=unpoller
AutoUpdate=registry
Pull=newer

# Configuration
Volume=%h/containers/config/unpoller/up.conf:/etc/unpoller/up.conf:Z,ro
Volume=%h/containers/data/unpoller:/data:Z

# Environment variables for UniFi controller credentials
# These should be set via podman secret or systemd EnvironmentFile
# For now, using placeholders - update before starting
Secret=unpoller_url,type=env,target=UP_UNIFI_DEFAULT_URL
Secret=unpoller_user,type=env,target=UP_UNIFI_DEFAULT_USER
Secret=unpoller_pass,type=env,target=UP_UNIFI_DEFAULT_PASS
Environment=TZ=Europe/Oslo

# Network configuration
Network=systemd-monitoring

# Prometheus metrics exposed on systemd-monitoring network (port 9130)
# No PublishPort needed - Prometheus can reach unpoller:9130 directly

# Health check disabled - distroless container lacks curl
# Monitoring via Prometheus scrape health instead
# HealthCmd=curl -f http://localhost:9130/metrics || exit 1
# HealthInterval=30s
# HealthRetries=3
# HealthStartPeriod=10s
# HealthTimeout=5s

# Security
User=1000:1000
NoNewPrivileges=true
SecurityLabelDisable=false

[Service]
# Resource limits
MemoryMax=256M
MemoryHigh=230M
CPUQuota=50%

# Restart policy
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=default.target
