# systemd/backup-restore-test.service
[Unit]
Description=Backup Restore Test
Documentation=file:///home/patriark/containers/docs/20-operations/guides/disaster-recovery.md
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/patriark/containers/scripts/test-backup-restore.sh --verbose
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2h

# Send alert on failure
# OnFailure=email-notification@%n.service

[Install]
WantedBy=default.target
