# systemd/backup-restore-test.service
[Unit]
Description=Backup Restore Test
Documentation=file:///home/patriark/containers/docs/20-operations/guides/disaster-recovery.md
After=network-online.target

[Service]
Type=oneshot
ExecStart=/home/patriark/containers/scripts/test-backup-restore.sh --verbose
StandardOutput=journal
StandardError=journal
TimeoutStartSec=2h

# Failure alerting is via Prometheus, NOT systemd OnFailure: the script emits
# backup_restore_test_success / _last_run_timestamp to the node_exporter textfile
# collector, and config/prometheus/alerts/backup-alerts.yml fires a critical alert on
# backup_restore_test_success==0 plus a 35-day staleness warning → alert-discord-relay.
# (email-notification@%n.service never existed — do not re-add it.)

[Install]
WantedBy=default.target
