# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare Dynamic DNS Update

[Service]
Type=oneshot
ExecStart=%h/containers/scripts/cloudflare-ddns.sh
