# systemd/db-dump.service
[Unit]
Description=Database application-consistent dumps (Tier 2 backup, ADR-029)
Documentation=file:///home/patriark/containers/docs/00-foundation/decisions
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/containers/scripts/db-dump.sh
# WP2 (L-019): coverage audit after each dump run — flags DB-like services
# missing from the ENGINES array via db_dump_coverage_gap metrics. "-" prefix:
# a gap alerts through Prometheus, it must not mark the dump run failed.
ExecStartPost=-%h/containers/scripts/audit-dump-coverage.sh
TimeoutStartSec=1h

# Low priority — do not interfere with interactive use or other I/O
Nice=19
IOSchedulingClass=idle

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
