# systemd/journal-logrotate.service
[Unit]
Description=Rotate journal export log file

[Service]
Type=oneshot
ExecStart=%h/containers/scripts/rotate-journal-export.sh
