# systemd/security-comprehensive-audit.timer
[Unit]
Description=Biweekly Comprehensive Security Audit

[Timer]
OnCalendar=*-*-01,15 06:45:00
Persistent=true

[Install]
WantedBy=timers.target
