# systemd/weekly-intelligence.timer
[Unit]
Description=Weekly Homelab Intelligence Report Timer
Documentation=file:///home/patriark/containers/scripts/weekly-intelligence-report.sh

[Timer]
# Run every Friday at 7:30 AM (start of workday, end of week summary)
OnCalendar=Fri *-*-* 07:30:00

# If system was off, run on next boot
Persistent=true

# Randomize by up to 5 minutes (keep predictable)
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
