# .claude/skills/homelab-deployment/templates/quadlets/background-worker.container
# ~/.config/containers/systemd/{{SERVICE_NAME}}.container
# Generated by homelab-deployment skill
# Background worker pattern (no web UI, internal processing)

# ═══════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════
[Unit]
Description={{SERVICE_DESCRIPTION}}
Documentation={{DOCS_URL}}
After=network-online.target {{NETWORK_SERVICES_AFTER}}
Wants=network-online.target
Requires={{NETWORK_SERVICES_REQUIRES}}
# Require dependencies (e.g., database, message queue)
{{REQUIRES_SERVICES}}

# ═══════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════
[Container]
Image={{IMAGE}}
ContainerName={{SERVICE_NAME}}
HostName={{SERVICE_NAME}}
AutoUpdate=registry

# Networks (no reverse_proxy - internal only)
# Static IP for consistent DNS resolution (ADR-018)
Network=systemd-{{PRIMARY_NETWORK}}:ip={{STATIC_IP_PRIMARY}}

# Volumes
Volume={{CONFIG_DIR}}:/config:Z
Volume={{DATA_DIR}}:/data:Z

# Environment variables
Environment=TZ=Europe/Oslo

# Health check (if applicable)
HealthCmd={{HEALTH_CMD}}
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=120s

# ═══════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════
[Service]
Slice=container.slice

# Resource limits
MemoryMax={{MEMORY_LIMIT}}
MemoryHigh={{MEMORY_HIGH}}

# Service behavior
Restart=on-failure
RestartSec=60s
TimeoutStartSec=300
TimeoutStopSec=120

# ═══════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════
[Install]
WantedBy=default.target
