# .claude/skills/homelab-deployment/templates/quadlets/database.container
# ~/.config/containers/systemd/{{SERVICE_NAME}}.container
# Generated by homelab-deployment skill
# Database service pattern (PostgreSQL, MySQL, Redis, etc.)

# ═══════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════
[Unit]
Description={{SERVICE_DESCRIPTION}}
Documentation={{DOCS_URL}}
After=network-online.target
Wants=network-online.target

# ═══════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════
[Container]
Image={{IMAGE}}
ContainerName={{SERVICE_NAME}}
HostName={{SERVICE_NAME}}
AutoUpdate=registry

# Networks (internal only - no reverse proxy access)
Network=systemd-database

# Volumes
# IMPORTANT: Database directories should have NOCOW attribute on BTRFS
# Run: chattr +C {{DATA_DIR}} before first use
Volume={{DATA_DIR}}:/var/lib/postgresql/data:Z

# Environment variables
Environment=TZ=Europe/Oslo
Environment=POSTGRES_DB={{DB_NAME}}
Environment=POSTGRES_USER={{DB_USER}}
Secret=db_password,type=env,target=POSTGRES_PASSWORD

# Health check
HealthCmd={{HEALTH_CMD}}
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# ═══════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════
[Service]
# Resource limits
MemoryMax={{MEMORY_LIMIT}}
MemoryHigh={{MEMORY_HIGH}}

# Service behavior
Restart=on-failure
RestartSec=30s
TimeoutStartSec=120
TimeoutStopSec=120

# ═══════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════
[Install]
WantedBy=default.target
