# .claude/skills/homelab-deployment/templates/quadlets/monitoring-service.container
# ~/.config/containers/systemd/{{SERVICE_NAME}}.container
# Generated by homelab-deployment skill
# Monitoring exporter/service pattern (Node Exporter, cAdvisor, etc.)

# ═══════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════
[Unit]
Description={{SERVICE_DESCRIPTION}}
Documentation={{DOCS_URL}}
After=network-online.target
Wants=network-online.target

# ═══════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════
[Container]
Image={{IMAGE}}
ContainerName={{SERVICE_NAME}}
HostName={{SERVICE_NAME}}
AutoUpdate=registry

# Networks (monitoring only - no external access)
Network=systemd-monitoring

# Volumes (minimal for exporters)
Volume={{CONFIG_DIR}}:/config:Z

# Environment variables
Environment=TZ=Europe/Oslo

# Metrics port (internal only)
PublishPort={{METRICS_PORT}}:{{METRICS_PORT}}

# Health check
HealthCmd={{HEALTH_CMD}}
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# ═══════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════
[Service]
# Resource limits (lightweight for exporters)
MemoryMax={{MEMORY_LIMIT}}
MemoryHigh={{MEMORY_HIGH}}

# Service behavior
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60
TimeoutStopSec=30

# ═══════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════
[Install]
WantedBy=default.target
