# .claude/skills/homelab-deployment/templates/quadlets/web-app.container
# ~/.config/containers/systemd/{{SERVICE_NAME}}.container
# Generated by homelab-deployment skill
# Standard web application pattern

# ═══════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════
[Unit]
Description={{SERVICE_DESCRIPTION}}
Documentation={{DOCS_URL}}
After=network-online.target
Wants=network-online.target

# ═══════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════
[Container]
Image={{IMAGE}}
ContainerName={{SERVICE_NAME}}
HostName={{SERVICE_NAME}}
AutoUpdate=registry

# Networks (first network provides default route)
Network=systemd-reverse_proxy
Network=systemd-monitoring

# Volumes
Volume={{CONFIG_DIR}}:/config:Z
Volume={{DATA_DIR}}:/data:Z

# Environment variables
Environment=TZ=Europe/Oslo

# Health check
HealthCmd={{HEALTH_CMD}}
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# ═══════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════
[Service]
# Resource limits
MemoryMax={{MEMORY_LIMIT}}
MemoryHigh={{MEMORY_HIGH}}

# Service behavior
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300
TimeoutStopSec=60

# ═══════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════
[Install]
WantedBy=default.target
