# .claude/skills/homelab-deployment/templates/quadlets/web-app.container
# ~/.config/containers/systemd/{{SERVICE_NAME}}.container
# Generated by homelab-deployment skill
# Standard web application pattern

# ═══════════════════════════════════════════════
# DEPENDENCIES
# ═══════════════════════════════════════════════
[Unit]
Description={{SERVICE_DESCRIPTION}}
Documentation={{DOCS_URL}}
After=network-online.target {{NETWORK_SERVICES_AFTER}}
Wants=network-online.target
Requires={{NETWORK_SERVICES_REQUIRES}}

# ═══════════════════════════════════════════════
# CONTAINER CONFIGURATION
# ═══════════════════════════════════════════════
[Container]
Image={{IMAGE}}
ContainerName={{SERVICE_NAME}}
HostName={{SERVICE_NAME}}
AutoUpdate=registry

# Networks (first network provides default route)
# Static IPs required for multi-network containers (ADR-018)
Network=systemd-reverse_proxy:ip={{STATIC_IP_REVERSE_PROXY}}
Network=systemd-monitoring:ip={{STATIC_IP_MONITORING}}

# Volumes
Volume={{CONFIG_DIR}}:/config:Z
Volume={{DATA_DIR}}:/data:Z

# Environment variables
Environment=TZ=Europe/Oslo

# Health check
HealthCmd={{HEALTH_CMD}}
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# ═══════════════════════════════════════════════
# SERVICE BEHAVIOR
# ═══════════════════════════════════════════════
[Service]
Slice=container.slice

# Resource limits
MemoryMax={{MEMORY_LIMIT}}
MemoryHigh={{MEMORY_HIGH}}
# Optional: MemorySwapMax={{MEMORY_SWAP_MAX}} (for resource-intensive services)

# Service behavior
Restart=on-failure
RestartSec=30s
TimeoutStartSec=300
TimeoutStopSec=60

# ═══════════════════════════════════════════════
# INSTALLATION
# ═══════════════════════════════════════════════
[Install]
WantedBy=default.target
