# ansible/roles/pihole_backup/templates/pihole-backup.service.j2
# Managed by Ansible (ADR-031 D8).
[Unit]
Description=Pi-hole Teleporter export + rsync push to Urd-backed host
Wants=network-online.target
After=network-online.target pihole-FTL.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pihole-backup.sh
# Surface failures: a missed backup should be visible (journal + systemd state).
